# quadlet/apps/hecate-traderw.container
[Unit]
Description=Hecate Trader Web Frontend (plugin)
After=hecate-traderd.service

[Container]
Image=ghcr.io/hecate-social/hecate-traderw:0.1.0
ContainerName=hecate-traderw
AutoUpdate=registry

# Web frontend serves on a port
PublishPort=5174:5174

# Mount traderd socket so nginx can proxy API calls
Volume=%h/.hecate/hecate-traderd/sockets:/var/run/traderd:ro

# Health check via HTTP
HealthCmd=curl -sf http://localhost:5174/ || exit 1
HealthInterval=30s
HealthRetries=3
HealthTimeout=5s
HealthStartPeriod=10s

[Service]
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s

[Install]
WantedBy=default.target
