# quadlet/apps/hecate-traderd.container
[Unit]
Description=Hecate Trader Daemon (plugin)
After=hecate-daemon.service
Wants=hecate-daemon.service

[Container]
Image=ghcr.io/hecate-social/hecate-traderd:647b54d
ContainerName=hecate-traderd
AutoUpdate=registry
Network=host

# Convention: mount to match $HOME/.hecate/hecate-traderd inside container
Volume=%h/.hecate/hecate-traderd:/app/.hecate/hecate-traderd:Z

# Daemon socket for plugin registration (convention path)
Volume=%h/.hecate/hecate-daemon/sockets:/app/.hecate/hecate-daemon/sockets:ro

# Health check via socket presence (convention path)
HealthCmd=test -S /app/.hecate/hecate-traderd/sockets/api.sock
HealthInterval=30s
HealthRetries=3
HealthTimeout=5s
HealthStartPeriod=10s

[Service]
Restart=on-failure
RestartSec=5s
TimeoutStartSec=60s

[Install]
WantedBy=default.target
