# quadlet/system/hecate-daemon.container
[Unit]
Description=Hecate Daemon (core)
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/hecate-social/hecate-daemon:0.9.2
ContainerName=hecate-daemon
AutoUpdate=registry
Network=host

# HOME=%h so convention paths (~/.hecate/...) resolve to real user home
Environment=HOME=%h

# Volume mount at real home path — convention paths identical inside/outside
Volume=%h/.hecate/hecate-daemon:%h/.hecate/hecate-daemon:Z

# Environment from files
EnvironmentFile=%h/.hecate/gitops/system/hecate-daemon.env
EnvironmentFile=-%h/.hecate/secrets/llm-providers.env

# Health check via convention path
HealthCmd=test -S %h/.hecate/hecate-daemon/sockets/api.sock
HealthInterval=30s
HealthRetries=3
HealthTimeout=5s
HealthStartPeriod=10s

[Service]
Restart=on-failure
RestartSec=5s
TimeoutStartSec=60s

[Install]
WantedBy=default.target
